# tap-broadcaster.service
[Unit]
Description=NFC Tap Broadcaster for Stuco POS
After=network.target

[Service]
Type=simple
User=qiss
WorkingDirectory=/home/qiss/stuco
Environment="NEXTJS_URL=http://localhost:3000"
Environment="NFC_TAP_SECRET=your-secret-here"
Environment="POS_LANE_ID=default"
Environment="PN532_DEVICE=i2c:/dev/i2c-1:pn532"
ExecStart=/home/qiss/stuco/.venv/bin/python -u /home/qiss/stuco/tap-broadcaster.py
StandardOutput=journal
StandardError=journal
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target

