# systemd/stuco-web.service
# ============================================================================
# INSTALLATION INSTRUCTIONS
# ============================================================================
# Before using this service file, you MUST customize the following:
# 1. Replace YOUR_USERNAME with your actual system username (e.g., pi, ubuntu)
# 2. Replace /path/to/stuco with your actual project installation path
#
# Example:
#   User=pi
#   WorkingDirectory=/home/pi/stuco/web-next
#   ExecStart=/usr/bin/node /home/pi/stuco/web-next/server.js
# ============================================================================

[Unit]
Description=SCPS Web Server (Next.js with WebSocket)
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/path/to/stuco/web-next
Environment="NODE_ENV=production"
Environment="PORT=3000"
Environment="HOSTNAME=0.0.0.0"
ExecStart=/usr/bin/node /path/to/stuco/web-next/server.js
StandardOutput=journal
StandardError=journal
Restart=on-failure
RestartSec=5

# Reload environment on restart
ExecStartPre=/bin/sleep 2

[Install]
WantedBy=multi-user.target

