# systemd/tap-broadcaster-reader2.service
# ============================================================================
# INSTALLATION INSTRUCTIONS
# ============================================================================
# Before using this service file, you MUST customize the following:
# 1. Replace YOUR_USERNAME with your actual system username (e.g., pi, ubuntu)
# 2. Replace /path/to/stuco with your actual project installation path
#
# Example:
#   User=pi
#   WorkingDirectory=/home/pi/stuco
#   ExecStart=/home/pi/stuco/.venv/bin/python -u /home/pi/stuco/tap-broadcaster.py
# ============================================================================

[Unit]
Description=NFC Tap Broadcaster Reader 2 for SCPS POS (WebSocket)
After=network.target

[Service]
Type=simple
User=qiss
WorkingDirectory=/home/qiss/scps
# Load environment variables from file (including NFC_TAP_SECRET)
EnvironmentFile=/home/qiss/scps/.env.broadcaster
Environment="NEXTJS_URL=http://localhost:3000"
Environment="POS_LANE_ID=reader-2"
Environment="PN532_DEVICE=tty:USB1:pn532"

# Pre-start: Wait for device to be ready (avoid race condition)
ExecStartPre=/bin/sleep 5

# Main process with unbuffered output
ExecStart=/home/qiss/scps/.venv/bin/python -u /home/qiss/scps/tap-broadcaster.py

# Post-stop: Allow cleanup time
ExecStopPost=/bin/sleep 1
StandardOutput=journal
StandardError=journal

# Restart policy: Always restart on failure
Restart=always
RestartSec=5

# Kill settings: Give time for cleanup
TimeoutStopSec=10
KillMode=mixed

[Install]
WantedBy=multi-user.target

