# systemd/tap-broadcaster-reader2.service
[Unit]
Description=NFC Tap Broadcaster Reader 2 for SCPS POS (WebSocket)
After=network.target

[Service]
Type=simple
User=qiss
WorkingDirectory=/home/qiss/stuco
# Load environment variables from file (including NFC_TAP_SECRET)
EnvironmentFile=/home/qiss/stuco/.env.broadcaster
Environment="NEXTJS_URL=http://localhost:3000"
Environment="POS_LANE_ID=reader-2"
Environment="PN532_DEVICE=tty:USB1:pn532"

# Pre-start: Reset USB device to ensure clean state
ExecStartPre=/home/qiss/stuco/scripts/reset-usb-nfc.sh USB1

# Main process with unbuffered output
ExecStart=/home/qiss/stuco/.venv/bin/python -u /home/qiss/stuco/tap-broadcaster.py

# Post-stop: Allow cleanup time
ExecStopPost=/bin/sleep 1
StandardOutput=journal
StandardError=journal

# Restart policy: Always restart on failure
Restart=always
RestartSec=5

# Kill settings: Give time for cleanup
TimeoutStopSec=10
KillMode=mixed

[Install]
WantedBy=multi-user.target

