# systemd/tap-broadcaster.service
# NOTE: Example service file - customize for your deployment before use
[Unit]
Description=NFC Tap Broadcaster for SCPS POS (WebSocket)
After=network.target

[Service]
Type=simple
User=stuco
WorkingDirectory=/opt/stuco
# Load environment variables from file (including NFC_TAP_SECRET)
EnvironmentFile=/opt/stuco/.env.broadcaster
Environment="NEXTJS_URL=http://localhost:3000"
Environment="POS_LANE_ID=reader-1"
Environment="PN532_DEVICE=tty:USB0:pn532"

# Pre-start: Reset USB device to ensure clean state
ExecStartPre=/opt/stuco/scripts/reset-usb-nfc.sh

# Main process with unbuffered output
ExecStart=/opt/stuco/.venv/bin/python -u /opt/stuco/tap-broadcaster.py

# Post-stop: Allow cleanup time
ExecStopPost=/bin/sleep 1
StandardOutput=journal
StandardError=journal

# Restart policy: Always restart on failure
Restart=always
RestartSec=5

# Kill settings: Give time for cleanup
TimeoutStopSec=10
KillMode=mixed

[Install]
WantedBy=multi-user.target

