# systemd/tap-broadcaster.service
[Unit]
Description=NFC Tap Broadcaster for Stuco POS (WebSocket)
After=network.target

[Service]
Type=simple
User=qiss
WorkingDirectory=/home/qiss/stuco
# Load environment variables from file (including NFC_TAP_SECRET)
EnvironmentFile=/home/qiss/stuco/.env.broadcaster
Environment="NEXTJS_URL=http://localhost:3000"
Environment="POS_LANE_ID=default"
Environment="PN532_DEVICE=tty:USB0:pn532"
ExecStart=/home/qiss/stuco/.venv/bin/python -u /home/qiss/stuco/tap-broadcaster.py
StandardOutput=journal
StandardError=journal
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target

